# deb/lib/systemd/system/home-measures.service.d/env.conf
[Service]
Environment="HM_DB_FILE=/var/lib/home-measures/hm.db"